# sentry-usbkey.rules
ACTION=="add", KERNEL=="sd*", RUN+="/usr/local/bin/sentry-usbkey udev"
