# sentry-usbkey.rules
ACTION=="add", KERNEL=="sd*", RUN+="/usr/lib/sentry-usbkey/sentry-usbkey-udev"
