# install/keyfile.rules.template
ACTION=="add", KERNEL=="sd*", RUN+="/PATH/TO/bin/udev-run-or-req"
